# launchd/services/com.jkorbmacher.vdirsyncer-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.jkorbmacher.vdirsyncer-sync</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/jkorbmacher/.config/launchd/scripts/vdirsyncer-sync</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>XDG_CACHE_DIR</key>
      <string>/Users/jkorbmacher/.cache</string>
      <key>XDG_CONFIG_HOME</key>
      <string>/Users/jkorbmacher/.config</string>
      <key>XDG_DATA_HOME</key>
      <string>/Users/jkorbmacher/.local/share</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>StandardOutPath</key>
    <string>/dev/null</string>
    <key>StandardErrorPath</key>
    <string>/dev/null</string>
  </dict>
</plist>
